# shajara.service
[Unit]
Description=Shajara — Django (Gunicorn)
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/shajara/backend
EnvironmentFile=/var/www/shajara/backend/.env
ExecStartPre=/var/www/shajara/venv/bin/python manage.py migrate --no-input
ExecStart=/var/www/shajara/venv/bin/gunicorn \
    config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
